# wg-x-daemon.service
[Unit]
Description=WireGuard X Daemon
Documentation=https://github.com/your-repo/wg-x
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wg-x daemon
Restart=on-failure
RestartSec=5s

# 安全设置
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run/wireguard /var/run

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
