# wg-x-daemon.service
[Unit]
Description=WireVault Daemon
Documentation=https://github.com/mrtian2016/wire-vault
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wire-vault daemon
Restart=on-failure
RestartSec=5s

# 安全设置
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run/wireguard /var/run

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
